# files/etc/systemd/system/music-backup.timer
[Unit]
Description=Daily backup of music library.

[Timer]
OnCalendar=*-*-* 01:35:00

[Install]
WantedBy=timers.target
